# argus.service
[Unit]
Description=Argus — Loki & Grafana Health Monitor
Documentation=https://github.com/siyamsarker/argus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=argus
Group=argus

WorkingDirectory=/opt/argus
ExecStart=/opt/argus/venv/bin/python /opt/argus/argus.py
EnvironmentFile=/opt/argus/.env

Restart=on-failure  # auto-restart if the process exits with a non-zero code
RestartSec=30       # wait 30 seconds before restarting

# Security hardening — restrict what the daemon process can do
NoNewPrivileges=true          # prevent privilege escalation via setuid binaries
ProtectSystem=strict          # mount / as read-only (except ReadWritePaths)
ProtectHome=true              # hide /home, /root, /run/user from the process
ReadWritePaths=/var/log/argus # allow writing to the log directory only
PrivateTmp=true               # isolate /tmp from other services

# Logging (stdout/stderr → journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argus

[Install]
WantedBy=multi-user.target
