# argus.service
[Unit]
Description=Argus — Loki & Grafana Health Monitor
Documentation=https://github.com/siyamsarker/argus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=argus
Group=argus

WorkingDirectory=/opt/argus
ExecStart=/opt/argus/venv/bin/python /opt/argus/argus.py
EnvironmentFile=/opt/argus/.env

# Auto-restart on non-zero exit; wait 30 seconds before retry
Restart=on-failure
RestartSec=30

# Security hardening — restrict what the daemon process can do
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Allow writing to the log directory only
ReadWritePaths=/var/log/argus

# Logging (stdout/stderr → journald)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argus

[Install]
WantedBy=multi-user.target
